# vuln-patch-agent_0.1.0/DEBIAN/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
    rm -f /etc/sudoers.d/patch-agent
fi

# disable systemd units on removal
if command -v systemctl >/dev/null 2>&1; then
    systemctl disable --now vuln-patch-agent.timer >/dev/null 2>&1 || true
    systemctl disable --now vuln-patch-agent.service >/dev/null 2>&1 || true
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

exit 0
