# vuln-patch-agent_0.1.0/lib/systemd/system/vuln-patch-agent.timer
[Unit]
Description=Run vuln-patch-agent daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=vuln-patch-agent.service

[Install]
WantedBy=timers.target
